# finance-tracker-sqlite-frontend.example.service
[Unit]
Description=Finance Tracker SQLite frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/user/WorkSpace/finance-tracker-sqlite
EnvironmentFile=/home/user/.config/finance-tracker-sqlite/systemd.env
ExecStart=/usr/bin/env bash -lc 'exec "$NPM_BIN" run sqlite:frontend -- --db "$DB_PATH" --user-id "$USER_ID" --bridge-host "$BRIDGE_HOST" --serve-host "$SERVE_HOST" --public-origin "$PUBLIC_ORIGIN" --login-email-env LOGIN_EMAIL --login-password-env LOGIN_PASSWORD'
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
